# config/.config/systemd/user/replace-wallpaper.service
[Unit]
Description=Replace Wallpaper

[Service]
Type=simple
ExecStart=/home/user/.bin/replacewallpaper

[Install]
WantedBy=default.target
