# config/.config/systemd/user/replace-wallpaper.service
[Unit]
Description=Replace Wallpaper

[Service]
Type=simple
ExecStart=/home/user/.bin/replacewallpaper

[Install]
After=network-online.target
Wants=network-online.target
